# deploy/strava-router.service
[Unit]
Description=Strava Heat Router (Vienna)
After=network.target

[Service]
Type=simple
User=%i
WorkingDirectory=/srv/strava-router
ExecStart=/srv/strava-router/.venv/bin/gunicorn -w 2 --threads 4 -b 127.0.0.1:8000 --chdir /srv/strava-router web.server:app
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/srv/strava-router
PrivateTmp=true

[Install]
WantedBy=multi-user.target
